# deploy/systemd/lafufu.target
[Unit]
Description=Lafufu (all services)
Requires=nats.service
After=nats.service

[Install]
WantedBy=multi-user.target
