# deploy/systemd/lafufu-animator.service
[Unit]
Description=Lafufu Animator (servos + lipsync)
# Order after control so request_config_snapshot fires AFTER control has
# subscribed to it on its own startup (otherwise the snapshot request is
# silently dropped — NATS is fire-and-forget — and the animator boots with
# env defaults, ignoring DB-tuned idle positions).
After=nats.service lafufu-control.service
Wants=nats.service
PartOf=lafufu.target
# StartLimit* must live in [Unit] on modern systemd; in [Service] they are
# silently ignored (the directives are parsed but have no effect).
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=lafufu
Group=lafufu
WorkingDirectory=/srv/lafufu
Environment=LAFUFU_NATS_URL=nats://localhost:4222
Environment=LAFUFU_DATA_DIR=/var/lafufu
ExecStart=/srv/lafufu/.venv/bin/python -m lafufu_animator
Restart=on-failure
RestartSec=5s
TimeoutStopSec=10s

[Install]
WantedBy=lafufu.target
