# deploy/systemd/lafufu-printer.service
[Unit]
Description=Lafufu Printer (CUPS bridge)
After=nats.service cups.service
Wants=nats.service
PartOf=lafufu.target

[Service]
Type=simple
User=lafufu
Group=lp
WorkingDirectory=/srv/lafufu
Environment=LAFUFU_NATS_URL=nats://localhost:4222
Environment=LAFUFU_PRINTER_AUTO=false
ExecStart=/srv/lafufu/.venv/bin/python -m lafufu_printer
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=lafufu.target
