# deploy/systemd/lafufu-agent.service
[Unit]
Description=Lafufu Agent (voice pipeline)
# Order after control so request_config_snapshot fires AFTER control has
# subscribed to it on its own startup (otherwise the snapshot request is
# silently dropped — NATS is fire-and-forget — and the agent boots with
# env defaults, ignoring DB-tuned settings).
After=nats.service lafufu-control.service
Wants=nats.service
PartOf=lafufu.target
# StartLimit* must live in [Unit] on modern systemd; in [Service] they are
# silently ignored (the directives are parsed but have no effect).
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=lafufu
Group=lafufu
WorkingDirectory=/srv/lafufu
Environment=LAFUFU_NATS_URL=nats://localhost:4222
Environment=LAFUFU_DATA_DIR=/var/lafufu
Environment=LAFUFU_WHISPER_MODEL=tiny.en
Environment=LAFUFU_STT_BACKEND=openai-whisper
Environment=LAFUFU_LLM_MODEL=qwen2.5:7b
Environment=LAFUFU_OLLAMA_URL=http://localhost:11434
Environment=LAFUFU_PIPER_MODEL=/srv/lafufu/models/lafufu_voice.onnx
ExecStart=/srv/lafufu/.venv/bin/python -m lafufu_agent
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s

[Install]
WantedBy=lafufu.target
