# deploy/systemd/lafufu-kiosk.service
[Unit]
Description=Lafufu Kiosk Browser (face view)
After=lafufu-control.service graphical-session.target
Wants=lafufu-control.service
PartOf=lafufu.target

[Service]
Type=simple
User=lafufu
Group=lafufu
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStartPre=/bin/sh -c 'until curl -sf http://localhost:8080/api/state/snapshot >/dev/null; do sleep 1; done'
# --kiosk + --app together = no chrome at all (no tab, no URL bar, no menu).
# Extra flags suppress every nudge / first-run prompt / accidental-input
# vector that could break the kiosk illusion.
ExecStart=/usr/bin/chromium \
  --kiosk \
  --start-fullscreen \
  --app=http://localhost:8080/face \
  --noerrdialogs \
  --disable-infobars \
  --disable-translate \
  --disable-features=TranslateUI,Translate \
  --disable-restore-session-state \
  --disable-session-crashed-bubble \
  --no-first-run \
  --no-default-browser-check \
  --password-store=basic \
  --overscroll-history-navigation=0 \
  --disable-pinch \
  --user-data-dir=/home/lafufu/.config/lafufu-kiosk
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=lafufu.target
